# deploy/zemer-overrides.timer
[Unit]
Description=Refresh zemer-search id-overrides several times a day (lightweight Firestore fetch, applied without restart)

[Timer]
Unit=zemer-overrides.service
# Every 10 min so a curated blockedContentIds change (e.g. hiding a women's playlist) takes effect within minutes.
# Shabbat is handled IN CODE: the service has ExecCondition=node harness/shabbat.mjs, which skips any tick from
# 20 min before candle lighting until havdalah using accurate weekly Brooklyn times (Hebcal) — so this can run
# every 10 min all week and simply no-ops through Shabbat.
OnCalendar=*-*-* *:0/10:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
