# deploy/zemer-playlists.service
[Unit]
Description=zemer-search community-playlist discovery + revalidate
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# --- EDIT WorkingDirectory for your deployment ---
WorkingDirectory=/opt/zemer-search
# CORPUS_DB defaults to <WorkingDirectory>/data/corpus.db (same default as the API unit); set it only if
# the DB lives elsewhere — it must match the API's.
# Environment=CORPUS_DB=/var/lib/zemer-search/corpus.db
# Discover from the topical seeds (data/playlist-seeds.json) + REVALIDATE all stored playlists (prune
# stale, refresh counts, regenerate the rejected-artists review report). Browse+search only — no /player —
# so this runs fine on a datacenter IP. Heavier artist-name sweeps are a manual op (SEEDS=both FIRSTNAMES=1).
Environment=SEEDS=topics
# N caps the seed count — set it ≥ the number of topical seeds (data/playlist-seeds.json) so none are
# truncated (the default is only 40). Bump if you add many more seed terms.
Environment=N=200
Environment=PAGES=2
Environment=REVALIDATE=1
# Re-fetch searches older than 24h so each scheduled run surfaces NEW playlists (otherwise forever-cached).
Environment=SEARCH_MAX_AGE_H=24
Environment=CONCURRENCY=5
Environment=MIN_INTERVAL_MS=200
Nice=10
IOSchedulingClass=idle
ExecCondition=/usr/bin/env node harness/shabbat.mjs
ExecStart=/usr/bin/env node harvester/playlists.mjs
TimeoutStartSec=4h
