# deploy/zemer-refresh-daily.timer
[Unit]
Description=Daily SHALLOW corpus refresh for zemer-search (fast, catches new releases)

[Timer]
Unit=zemer-refresh@shallow.service
# Daily shallow pass Mon–Sat 03:00 (Sunday is the separate deep pass; both share a flock). Shabbat is handled
# IN CODE, not here: each maintenance service has ExecCondition=node harness/shabbat.mjs, which skips the run
# from 20 min before candle lighting until havdalah using accurate weekly Brooklyn times (Hebcal). So the Sat
# 03:00 tick that lands in Shabbat is skipped automatically, and Persistent catch-ups are gated the same way.
OnCalendar=Mon..Sat *-*-* 03:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
